# setting_files/etc/alias/profile.d/alias.sh
alias sudo='sudo '
alias cp='cp --reflink=auto --sparse=always'
